# systemd/backup.vps.timer
[Unit]
Description=Run backup routine weekly at 10:00 AM on Sunday

[Timer]
OnCalendar=Sun 10:00
Persistent=true

[Install]
WantedBy=timers.target
