# systemd/backup.pc.timer
[Unit]
Description=Run backup routine daily at 12:30 AM

[Timer]
OnCalendar=00:30
Persistent=true

[Install]
WantedBy=timers.target
