# systemd/backup.vps.timer
[Unit]
Description=Run backup routine daily at 5:00 AM PST

[Timer]
OnCalendar=13:00
Persistent=true

[Install]
WantedBy=timers.target
